# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/var/www/app/console server:run 0.0.0.0:8000